# contrib/service/launchd/com.forgente.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
		<key>Label</key>
		<string>com.forgente.web</string>
		<!-- assumes Forgente is running under 'git' account -->
		<!-- modify below to reflect your settings -->
		<key>UserName</key>
		<string>git</string>
		<key>GroupName</key>
		<string>git</string>
		<key>ProgramArguments</key>
		<array>
			<!-- assumes Forgente is installed in /Users/git/forgente -->
			<!-- modify below to reflect your settings -->
			<string>/Users/git/forgente/forgente</string>
			<string>web</string>
		</array>
		<key>RunAtLoad</key>
		<true/>
		<key>KeepAlive</key>
		<true/>
		<!-- assumes Forgente is installed in /Users/git/forgente -->
		<!-- modify below to reflect your settings -->
		<key>WorkingDirectory</key>
		<string>/Users/git/forgente/</string>
		<key>StandardOutPath</key>
		<string>/Users/git/forgente/log/stdout.log</string>
		<key>StandardErrorPath</key>
		<string>/Users/git/forgente/log/stderr.log</string>
		<!-- default 256 is too low for Forgente needs using parallel pipes -->
		<key>SoftResourceLimits</key>
		<dict>
			<key>NumberOfFiles</key>
			<integer>8192</integer>
		</dict>
	</dict>
</plist>
