# docker/root/etc/s6/forgente/run
#!/bin/bash
[[ -f ./setup ]] && source ./setup

pushd /app/forgente >/dev/null
exec su-exec $USER /usr/local/bin/forgente web
popd
